# solr/supervisor.d/init.ini
[program:init]
user = solr
command = /docker-entrypoint-initdb.d/init_solr
autorestart = false
autostart = true