# solr/supervisor.d/solr.ini
[program:solr]
user = solr
autostart = true
autorestart = true
command = SOLR_INCLUDE=/etc/default/solr.in.sh /opt/solr/bin/solr start