# solr/supervisor.d/solr.ini
[program:solr]
user = solr
autostart = true
autorestart = true
command = /opt/docker-solr/scripts/solr-foreground